# init/nfprom.openrc
#!/sbin/openrc-run

description="Prometheus IPTables stats exporter"

pidfile="/run/${SVCNAME}.pid"
command=/usr/sbin/nfprom
command_args="${NFPROM_ARGS}"
command_background="true"

start_stop_daemon_args="-u ${PROMETHEUS_USER} -1 ${NFPROM_LOG_FILE} -2 ${NFPROM_LOG_FILE} -- ${NFPROM_ARGS}"

